# .platform/files/supervisor.ini
[program:laravel-queue]
process_name=%(program_name)s
command=/usr/bin/php /var/www/html/artisan schedule:run 1>> /dev/null 2>&1
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stopwaitsecs=3600
user=root
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/worker.log