# .platform/files/supervisor.ini
[program:laravel-queue]
process_name=%(program_name)s
command=nice -n 10 /usr/bin/php /var/www/html/artisan queue:work  --sleep=3 --tries=2  --timeout=360 --rest=1 --max-jobs=1000 --max-time=3600 
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stopwaitsecs=3600
user=root
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/worker.log
